# deploy/systemd/elitev2-shorts-cleanup.timer
[Unit]
Description=Run elitev2-shorts-cleanup daily

[Timer]
OnCalendar=*-*-* 05:00:00
RandomizedDelaySec=10min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
