# deploy/systemd/elitev2-shorts-poll.timer
[Unit]
Description=Run elitev2-shorts-poll every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
