# deploy/systemd/elitev2-posts-dupescan.service
[Unit]
Description=Scan the elite-v2 posts library for duplicate images
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=thomas
Group=thomas
# Runs inside the container (shares the /posts-store mount + better-sqlite3).
# Groups byte-identical and perceptually identical images into post_dupe_groups
# for admin review under the posts Settings -> Duplicates. Deletes nothing. The
# per-image fingerprint cache makes every run after the first cheap (only
# new/changed images are hashed + perceptual-hashed).
ExecStart=/usr/bin/docker exec elitev2 node scripts/scan-posts-duplicates.mjs
Nice=18
IOSchedulingClass=idle
